# systemd/libsrc-mcp.service
[Unit]
Description=libsrc-mcp - Library Source Code MCP Server

[Service]
Type=simple
ExecStart=%h/.local/bin/uvx libsrc-mcp serve
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
